# crontab
0 * * * * root cd /backend && /usr/local/bin/python manage.py run_timed_tasks >> /var/log/cron.log 2>&1
